# files/frontend/vite.service.jinja
[Unit]
After=network.target

[Service]
User={{ user }}
Group=www-data
WorkingDirectory={{ frontend_workspace }}
ExecStart={{ frontend_exec_start }}

[Install]
WantedBy=multi-user.target