# files/backend/uvicorn.service.jinja
[Unit]
After=network.target

[Service]
User={{ user }}
Group=www-data
WorkingDirectory={{ workspace }}
ExecStart={{ exec_start }}

[Install]
WantedBy=multi-user.target